# ansible/templates/per-docker-roles/per-docker-roles.service.j2
[Unit]
Description=Per Docker Roles Credential Helper
Documentation=https://github.com/thepatrick/per-docker-roles
After=network.target

[Service]
Type=simple
User={{ user }}
ExecStart={{ binary }} serve --listen-on {{ listen_on }} --docker-network {{ docker_network }}
Restart=on-failure

[Install]
WantedBy=multi-user.target
